# config/bankease.socket
[Unit]
Description=BankEase gunicorn socket

[Socket]
ListenStream=/run/bankease.sock

[Install]
WantedBy=sockets.target
