# config/gunicorn.socket
[Unit]
Description=BankEase gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
